# systemd/display-rotate.timer
[Unit]
Description=Run display rotation 1 minute after boot

[Timer]
OnBootSec=1min
Persistent=false

[Install]
WantedBy=timers.target
